# deploy/portfolio-webhook.service
# =============================================================================
# portfolio-webhook.service — GitHub webhook listener (adnanh/webhook)
#
# Este archivo lo procesa install_service.sh, que sustituye:
#   __PROJECT_DIR__  →  ruta absoluta al proyecto
#   __USER__         →  usuario que ejecuta el servicio
# =============================================================================
[Unit]
Description=Portfolio Tracker — GitHub Webhook Listener
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# El secreto vive fuera de git (webhook.env). El '-' hace opcional el fichero
# para que el fallo, si falta, lo dé render_hooks.sh con un mensaje claro.
EnvironmentFile=-__PROJECT_DIR__/deploy/webhook.env
# Regenera hooks.json (git-ignorado) con el secreto real ANTES de arrancar.
# Así sobrevive a `git reset --hard` de update.sh: se reconstruye en cada start.
ExecStartPre=__PROJECT_DIR__/deploy/render_hooks.sh
# -urlprefix "" sirve el hook en /deploy-finance (no /hooks/deploy-finance):
# el Tailscale Funnel monta /hooks → :9000 quitando el prefijo /hooks, así que
# GitHub apunta a https://<host>/hooks/deploy-finance y, tras el strip, llega
# como /deploy-finance. Sin esto, el hook nunca casa y el auto-deploy no corre.
ExecStart=/usr/bin/webhook -hooks __PROJECT_DIR__/deploy/hooks.json -urlprefix "" -port 9000 -verbose
WorkingDirectory=__PROJECT_DIR__
User=__USER__
Environment=PATH=/usr/local/bin:/usr/bin:/bin

Restart=on-failure
RestartSec=10

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
